# buildroot-external/rootfs-overlay/usr/lib/systemd/system/check-hassos-supervisor.service
[Unit]
Description=Check that hassos-supervisor.service is alive

[Service]
Type=oneshot
ExecStart=/usr/bin/check-supervisor.sh
