# buildroot-external/rootfs-overlay/usr/lib/systemd/system/serial-getty@.service.d/haos.conf
[Unit]
Wants=haos-supervisor.service
After=haos-supervisor.service
